# tools/alpine/resizemmc
#!/sbin/openrc-run

description="Resize rootfs partition."

start() {
	if [ -e /root/.need_resize ] ; then
		ebegin "Resize rootfs, please reboot after finish"
		/sbin/resizemmc.sh > /dev/null 2>&1
		eend $?
		rm -f /root/.need_resize
	fi
}
